# meta-mct/meta-s8053/recipes-core/dropbear/dropbear/dropbear.socket
[Unit]
Conflicts=dropbear.service

[Socket]
ListenStream=22
Accept=yes
BindToDevice=eth0

[Install]
WantedBy=sockets.target
Also=dropbearkey.service
