# meta-mct/meta-s8033/recipes-phosphor/mct-default-action/mct-default-action/setting-hsc-register@.service
[Unit]
Description=%I hot swap controller via register

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/env setting-hsc-register.sh %I