# meta-mct/meta-s8033/recipes-phosphor/sensors/dbus-sensors/service_files/xyz.openbmc_project.eventsensor.service
[Unit]
Description=IPMI Event-only Sensor
StopWhenUnneeded=false
After=xyz.openbmc_project.EntityManager.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/eventsensor

[Install]
WantedBy=multi-user.target
