# meta-mct/meta-s8056/recipes-mct/default-action/mct-default-action/usb-init.service
[Unit]
Description=Host inventory service
StopWhenUnneeded=false
Requires=xyz.openbmc_project.EntityManager.service
After=xyz.openbmc_project.EntityManager.service

[Service]
Restart=simple
ExecStart=sh /usr/sbin/configure_usb_gadget.sh

[Install]
WantedBy=multi-user.target