# meta-mct/meta-s8056/recipes-mct/gpio/mct-gpio-monitor/gpio-prochot-thermtrip-deassert@.service
[Unit]
Description=Deassert %I 
StartLimitIntervalSec=6
StartLimitBurst=1

[Service]
Type=oneshot
ExecStart=/usr/bin/env prochotThermtripHandler deassert CPU_STATE %I
SyslogIdentifier=gpio-prochot-thermtrip
RemainAfterExit=true