# meta-mct/meta-s8056/recipes-phosphor/configuration/entity-manager/service_files/xyz.openbmc_project.FruDevice.service
[Unit]
Description=Fru Device
After=xyz.openbmc_project.CloseMuxes.service dbus.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/fru-device

[Install]
WantedBy=multi-user.target
