# meta-mct/meta-s8033/recipes-phosphor/mct-power-handler/mct-power-handler/obmc-chassis-poweroff.target
[Unit]
Description=Chassis (Power Off)
RefuseManualStart=no
RefuseManualStop=yes