# meta-mct/meta-s8033/recipes-phosphor/mct-power-handler/mct-power-handler/obmc-host-poweron.target
[Unit]
Description=Host (Power On)
RefuseManualStart=no
RefuseManualStop=yes