# meta-mct/meta-s8033/recipes-phosphor/mct-power-handler/mct-power-handler/obmc-host-poweroff.target
[Unit]
Description=Host (Power Off)
RefuseManualStart=no
RefuseManualStop=yes