# meta-mct/meta-s8033/recipes-phosphor/mct-power-handler/mct-power-handler/obmc-chassis-poweron.target
[Unit]
Description=Chassis (Power On)
RefuseManualStart=no
RefuseManualStop=yes