# meta-mct/meta-s8053/recipes-phosphor/fans/phosphor-pid-control/setting-fan-init.service
[Unit]
Description=Setting fan duty when PID service initializing

[Service]
Type=simple
ExecStart=/usr/bin/env controlFan.sh service-init
SyslogIdentifier=setting-fan-init

[Install]
WantedBy=phosphor-pid-control.service