# meta-mct/meta-s8033/recipes-phosphor/gpio/mct-gpio-monitor/gpio-prochot-thermtrip-assert@.service
[Unit]
Description=Assert %I
StartLimitIntervalSec=6
StartLimitBurst=1

[Service]
Type=oneshot
ExecStart=/usr/bin/env prochotThermtripHandler assert CPU_State %I
SyslogIdentifier=gpio-prochot-thermtrip
RemainAfterExit=true