# meta-mct/meta-s8053/recipes-mct/flash/mct-bmc-update-in-band/phosphor-ipmi-flash-alt-bmc-verify.service
[Unit]
Description=Phosphor-ipmi-flash bmc verify service

[Service]
ExecStart=/usr/bin/image-verify.sh "/run/initramfs/bmc-backup"
Type=oneshot

[Install]
WantedBy=phosphor-ipmi-flash-bmc-verify.target
